# systemd/soul-courier.service
[Unit]
Description=Soul Team Courier — agent message bus daemon
Documentation=https://github.com/your-org/soul-team
After=network.target
Wants=network.target

[Service]
Type=simple
User=%i
# Edit this path to match your installation
WorkingDirectory=%h/.local/share/soul-team
ExecStart=/usr/bin/python3 %h/.local/share/soul-team/courier/daemon.py \
    --config %h/.config/soul-team/cluster.yaml \
    --log-level info
Restart=on-failure
RestartSec=5s
# Prevent runaway restarts
StartLimitIntervalSec=60
StartLimitBurst=5

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=%h/.local/share/soul-team
ReadWritePaths=%h/.config/soul-team

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target

# Usage:
#   systemctl --user enable soul-courier
#   systemctl --user start soul-courier
#   systemctl --user status soul-courier
#   journalctl --user -u soul-courier -f
#
# Or run setup.sh which handles all of this automatically.
