# systemd/soul-guardian.service
[Unit]
Description=Soul Team Guardian — agent health watchdog
Documentation=https://github.com/rishav1305/soul-team
After=network.target soul-courier.service
Wants=soul-courier.service

[Service]
Type=simple
WorkingDirectory=%h/soul-team
ExecStart=/usr/bin/python3 %h/soul-team/guardian/guardian.py
Restart=on-failure
RestartSec=10s
StartLimitIntervalSec=120
StartLimitBurst=3

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=%h/soul-team
ReadWritePaths=%h/.claude/config
ReadWritePaths=%h/.soul-team
ReadWritePaths=%h/.soul
ReadWritePaths=%h/.clawteam

# Allow reading system temperature files (for hardware monitoring)
ReadOnlyPaths=/sys/class/thermal
ReadOnlyPaths=/proc

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
