# systemd/soul-courier.service
[Unit]
Description=Soul Team Courier — agent message bus daemon
Documentation=https://github.com/rishav1305/soul-team
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=%h/soul-team
ExecStart=/usr/bin/python3 %h/soul-team/courier/main.py
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60
StartLimitBurst=5

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=%h/soul-team
ReadWritePaths=%h/.claude/config
ReadWritePaths=%h/.soul
ReadWritePaths=%h/.clawteam

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
